# internal/provisioner/maintenance/systemd/provision-success.service
[Unit]
Description=Shoal Provisioner: Report Successful Provisioning
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=10min
StartLimitBurst=10

[Service]
Type=oneshot
EnvironmentFile=/run/provision/recipe.env
ExecStart=/usr/bin/env bash /opt/shoal/bin/send-webhook.sh success
Restart=on-failure
RestartSec=10s
RemainAfterExit=no
TimeoutSec=2min

[Install]
WantedBy=install-linux.target
